# web-build/start-vnc.conf
[program:start-vnc]
command=/usr/local/bin/start-vnc.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/start-vnc.log
stderr_logfile=/var/log/supervisor/start-vnc_err.log
